# ranch_os/rootfs-overlay/etc/systemd/system/systemd-timesyncd.service.d/ranch-os.conf
# Ranch OS override for systemd-timesyncd on Buildroot.
#
# THE BUG IT FIXES
# ----------------
# Buildroot's default rootfs ships /etc/resolv.conf as a symlink to
# /tmp/resolv.conf (so the read-only /etc can host a "dynamic" file).
# NetworkManager populates the symlink target with real nameservers
# at runtime.
#
# timesyncd's upstream unit file sets PrivateTmp=yes for sandboxing.
# That gives the service its own /tmp namespace — empty — and the
# /etc/resolv.conf symlink now points at a file that doesn't exist
# inside the sandbox. NSS lookups fail with "Temporary failure in
# name resolution", every hostname resolution attempt fails, and
# timesyncd never contacts an NTP server.
#
# Visible symptom (with SYSTEMD_LOG_LEVEL=debug):
#   timesyncd[N]: Resolving time1.google.com...
#   timesyncd[N]: Failed to resolve time1.google.com:
#                 Temporary failure in name resolution
#
# THE FIX
# -------
# Disable PrivateTmp for this one service so the resolv.conf symlink
# resolves correctly. timesyncd doesn't write anything sensitive to
# /tmp, so losing this hardening on it specifically is fine.
#
# The proper-upstream fix would be: don't ship a symlink-into-/tmp,
# or have NetworkManager write /etc/resolv.conf as a real file (which
# rc-manager=file *should* do, except it follows the existing symlink
# and writes through it). Both belong in Buildroot, not here.

[Service]
PrivateTmp=no
